# launchd/local.mcp-1password.plist.example
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>local.mcp-1password</string>

  <key>ProgramArguments</key>
  <array>
    <string>mcp-1password</string>
    <string>--transport=http</string>
    <string>--http-host=127.0.0.1</string>
    <string>--http-port=17337</string>
    <string>--auth-mode=desktop</string>
    <string>--account=YOUR_ACCOUNT_ID_OR_UUID</string>
    <string>--enable-script-runner=true</string>
    <string>--script-runner-root=/absolute/path/to/trusted/project/root</string>
    <string>--script-runner-allowlist=/absolute/path/to/trusted/project/root/.onepassword-mcp.json</string>
    <string>--script-runner-root=/absolute/path/to/another/trusted/project</string>
    <string>--script-runner-allowlist=/absolute/path/to/another/trusted/project/.onepassword-mcp.json</string>
    <string>--http-path=/mcp</string>
    <string>--op-cli-path=/absolute/path/to/op</string>
    <string>--op-cli-auth-mode=desktop</string>
  </array>

  <key>EnvironmentVariables</key>
  <dict>
    <key>HOME</key>
    <string>/absolute/path/to/home</string>
    <key>OP_MCP_HTTP_BEARER_TOKEN</key>
    <string>YOUR_LONG_RANDOM_TOKEN</string>
    <key>PATH</key>
    <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
  </dict>

  <key>WorkingDirectory</key>
  <string>/absolute/path/to/mcp-1password</string>

  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/absolute/path/to/Library/Logs/mcp-1password.log</string>
  <key>StandardErrorPath</key>
  <string>/absolute/path/to/Library/Logs/mcp-1password.err.log</string>
</dict>
</plist>
